# docker/config/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[include]
files = /etc/supervisor.d/*.conf

[program:fortress-panel-backend]
command=/usr/bin/node /app/dist/index.js
directory=/app
user=node
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/fortress-panel/backend.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=NODE_ENV=production

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/fortress-panel/nginx.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

[program:php-fpm]
command=/usr/sbin/php-fpm8 --nodaemonize --fpm-config /etc/php8/php-fpm.conf
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/fortress-panel/php-fpm.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

[program:redis]
command=redis-server --daemonize no --port 6379 --bind 127.0.0.1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/fortress-panel/redis.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

[program:docker-daemon]
command=/usr/bin/dockerd --host=unix:///var/run/docker.sock --host=tcp://0.0.0.0:2375 --storage-driver=overlay2
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/fortress-panel/docker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
user=root

[group:fortress-panel]
programs=fortress-panel-backend,nginx,php-fpm,redis,docker-daemon